# .docker/crontab
* * * * * /usr/local/bin/php /var/www/html/cron.php  >> /var/www/html/logs/cron.log 2>&1
